# 99-hp-wmi-permissions.rules
# /etc/udev/rules.d/99-hp-wmi-permissions.rules
#
# Grant victus group write access to HP WMI fan control and keyboard backlight.
# Uses /usr/bin/ explicit paths — required on Fedora where udev's minimal
# environment does not reliably resolve /bin symlinks to /usr/bin.

# Grant access to HP WMI hwmon devices to the victus group
SUBSYSTEM=="hwmon", KERNELS=="hp-wmi", GROUP="victus", MODE="0664"
SUBSYSTEM=="hwmon", KERNELS=="hp-wmi", ACTION=="add|change", RUN+="/usr/bin/chgrp -R victus /sys%p", RUN+="/usr/bin/chmod g+w /sys%p/fan1_target /sys%p/fan2_target"

# Grant access to HP keyboard backlight to the victus group
SUBSYSTEM=="leds", KERNEL=="hp::kbd_backlight", RUN+="/usr/bin/chgrp victus /sys%p/brightness /sys%p/multi_intensity", RUN+="/usr/bin/chmod g+w /sys%p/brightness /sys%p/multi_intensity"
